# shared/dot-config/systemd/user/xmonad@.target
[Unit]
Description=xmonad DISPLAY=%i
RefuseManualStart=yes
StopWhenUnneeded=yes

Wants=xsettingsd@%i.service
Wants=xss-lock@%i.service

