# linux/dot-config/systemd/user/juank-media-mpris-daemon.service
[Unit]
Description=juank-media mpris-daemon

[Service]
ExecStart=%h/bin/juank-media mpris-daemon
Restart=always

[Install]
WantedBy=default.target
