# dot-config/systemd/user/battery-watch.service
[Unit]
Description=battery-watch
PartOf=graphical-session.target
Wants=dunst.service
After=dunst.service

[Service]
ExecStart=%h/bin/battery-watch
Restart=always

[Install]
WantedBy=graphical-session.target
